# templates/sonarr.service.j2
[Service]
ExecStart=
ExecStart=/usr/bin/sonarr -nobrowser -data={{ sonarr_data_dir | quote }}
